# etc/pulp_workers.service.d/pulp.conf
[Unit]
Requires=var-lib-pulp.mount
After=var-lib-pulp.mount
