# etc/mongod.service.d/var-lib-mongodb.mount.conf
[Unit]
Wants=var-lib-mongodb.mount

[Service]
ExecStartPre=/usr/bin/sudo /sbin/blockdev --setra 8 /dev/disk/by-partlabel/mongodb
